# data/citadel-realmsd.service
[Unit]
Description=Realms Manager Daemon

[Service]
ExecStart=/usr/libexec/citadel-realmsd

[Install]
WantedBy=multi-user.target
